# supervisord.conf
[supervisord]
nodaemon=true

[program:gqa]
command=/go/bin/GQA
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:marvin]
command=/marvin/marvin
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0